# scripts/com.nateprich.fantasy.snapshot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nateprich.fantasy.snapshot</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/nateprich/code/Fantasy-Football/scripts/weekly-snapshot.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/nateprich/code/Fantasy-Football/logs/snapshot.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/nateprich/code/Fantasy-Football/logs/snapshot.stderr.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/nateprich/code/Fantasy-Football</string>
</dict>
</plist>
